# VKR/7.Ansible_automation/roles/sysvol_replication/templates/sysvol-sync.timer.j2
[Unit]
Description=Cинхронизация sysvol каждые 5 минуты с основного DC

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=sysvol-sync.service

[Install]
WantedBy=timers.target
